# openbmc/build/tmp/deploy/sdk/witherspoon-2019-08-08/sysroots/armv6-openbmc-linux-gnueabi/lib/systemd/system/check-pnor-format.service
[Unit]
Description=Check PNOR UBI format
Before=mboxd.service

[Service]
ExecStart=/usr/bin/check_pnor_format.sh
Type=oneshot
RemainAfterExit=yes

[Install]
RequiredBy=mboxd.service
