# openbmc/build/tmp/deploy/sdk/witherspoon-2019-08-08/sysroots/armv6-openbmc-linux-gnueabi/lib/systemd/system/obmc-flash-bmc-ubirw.service
[Unit]
Description=Create BMC read-write ubi volume

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/obmc-flash-bmc ubirw bmc rwfs 0x600000
