# openbmc/build/tmp/deploy/sdk/witherspoon-2019-08-08/sysroots/armv6-openbmc-linux-gnueabi/lib/systemd/system/org.openbmc.HostIpmi.service
[Unit]
Description=Phosphor IPMI BT DBus Bridge

[Service]
Restart=always
ExecStart=/usr/bin/env btbridged
SyslogIdentifier=btbridged
Type=dbus
BusName=org.openbmc.HostIpmi

[Install]
WantedBy=multi-user.target
