# openbmc/build/tmp/deploy/sdk/witherspoon-2019-08-08/sysroots/armv6-openbmc-linux-gnueabi/lib/systemd/system/xyz.openbmc_project.User.Manager.service
[Unit]
Description=Phosphor User Manager

[Service]
ExecStart=/usr/bin/env phosphor-user-manager
SyslogIdentifier=phosphor-user-manager
Restart=always
Type=dbus
BusName=xyz.openbmc_project.User.Manager

[Install]
WantedBy=multi-user.target
