# openbmc/build/tmp/deploy/sdk/witherspoon-2019-08-08/sysroots/armv6-openbmc-linux-gnueabi/lib/systemd/system/xyz.openbmc_project.Network.service
[Unit]
Description=Phosphor Network Manager
After=obmc-mapper.target

[Service]
ExecStart=/usr/bin/env phosphor-network-manager
SyslogIdentifier=phosphor-network-manager
Restart=always
Type=dbus
BusName=xyz.openbmc_project.Network
RuntimeDirectory = network
RuntimeDirectoryPreserve = yes
StateDirectory = network

[Install]
WantedBy=multi-user.target
