# openbmc/build/tmp/deploy/sdk/witherspoon-2019-08-08/sysroots/armv6-openbmc-linux-gnueabi/lib/systemd/system/xyz.openbmc_project.Software.Download.service
[Unit]
Description=Phosphor Download Manager

[Service]
ExecStart=/usr/bin/phosphor-download-manager
Restart=always
Type=dbus
BusName=xyz.openbmc_project.Software.Download

[Install]
WantedBy=multi-user.target
