# openbmc/build/tmp/deploy/sdk/witherspoon-2019-08-08/sysroots/armv6-openbmc-linux-gnueabi/lib/systemd/system/xyz.openbmc_project.LDAP.PrivilegeMapper.service
[Unit]
Description=Phosphor LDAP privilege mapper

[Service]
ExecStart=/usr/bin/env phosphor-ldap-mapper
SyslogIdentifier=phosphor-ldap-mapper
Restart=always
Type=dbus
BusName=xyz.openbmc_project.LDAP.PrivilegeMapper

[Install]
WantedBy=multi-user.target
